# ec2arch-systemd/root-passwd.service
[Unit]
Description=Randomize root password on first boot

[Service]
ExecStart=/usr/bin/randomize-root-passwd

[Install]
WantedBy=multi-user.target
